# ScoringEngine/DEBIAN/prerm
#!/bin/bash
# Author: Bailey Kasin

systemctl stop score_web
systemctl stop scoring
systemctl disable score_web
systemctl disable scoring